# web-framework/php/wordpress/src/code/root/lib/systemd/system/phpsessionclean.timer
[Unit]
Description=Clean PHP session files every 30 mins

[Timer]
OnCalendar=*-*-* *:09,39:00
Persistent=true

[Install]
WantedBy=timers.target
